# .config/systemd/user/boinccmd-run@.service
[Unit]
Description=boinccmd-run-%i

[Service]
Type=oneshot
ExecStart=/usr/bin/boinccmd --set_run_mode %i
