# .config/systemd/user/boinccmd-network@.service
[Unit]
Description=boinccmd-network-%i

[Service]
Type=oneshot
ExecStart=/usr/bin/boinccmd --set_network_mode %i
