# .config/systemd/user/deluged.service
[Unit]
Description=dfanjul deluged

[Service]
ExecStart=/usr/bin/deluged -d
LimitNOFILE=65536
