# .config/systemd/user/moosic.service
[Unit]
Description=dfanjul moosic

[Service]
ExecStart=/usr/bin/moosic start-server -f
