# .config/systemd/user/moosic.service
[Unit]
Description=moosic

[Service]
ExecStart=/usr/bin/moosic start-server -f
