# etc/cron.hourly/00-check-user-env
#!/bin/sh

[ "$USER" ] || echo '$USER is not set in cron'
