# etc/cron.hourly/00-check-user-env
#!/bin/bash
set -e

[ "$USER" ] || echo '$USER is not set in cron'
