# etc/cron.hourly/01-newsbeuter
#!/bin/bash
set -e

mkdir -p ~/var/{lock,log}
lock=newsbeuter-reload
nohup postpone -L ~/var/lock/elock-"$lock" newsbeuter -x reload &>> ~/var/log/"$lock".log &
