# etc/cron.hourly/01-newsbeuter
#!/bin/bash

lock=newsbeuter-reload
nohup postpone -l ~/var/lock/lock-"$lock" newsbeuter -x reload &>> ~/var/log/"$lock".log &
