# etc/cron.hourly/01-newsbeuter
#!/bin/sh

tsp -m newsbeuter -x reload >/dev/null
