# etc/cron.hourly/50-backup-dir-var-backups-home
#!/bin/bash
set -e

[ -d /home/dfanjul/var/backups/home ] || exit 0
~/bin/backup-dir-postpone \
  ~ file:///home/dfanjul/var/backups/home \
  --volsize 64 --full-if-older-than 13D --exclude-globbing-filelist ~/etc/exclude.home.list \
  2>/dev/null &
