# etc/cron.hourly/80-git-fsck-and-gc
#!/bin/bash
set -e

TERM= xargs -rn 1 -I folder bash -c \
  'cd "folder" &>/dev/null && ~/bin/git-fsck-and-gc' \
  < ~/var/git-fsck-and-gc.list
