# etc/cron.hourly/bundle-backup-home
#!/bin/bash

tsp -m ~/bin/bundle-backup-home >/dev/null
