# etc/cron.hourly/bup-backup-home
#!/bin/bash

tsp -m ~/bin/bup-backup-home >/dev/null
