# etc/cron.hourly/git-fsck-and-gc-all
#!/bin/bash

tsp -m ~/bin/git-fsck-and-gc-all >/dev/null
