# etc/cron.hourly/rsync-backup-home
#!/bin/bash

TERM="" tsp -m ~/bin/rsync-backup-home >/dev/null
