# etc/cron.hourly/rsync-backup-home
#!/bin/bash

tsp -m ~/bin/rsync-backup-home-all >/dev/null
