# etc/cron.hourly/uprecords
#!/bin/sh

uprecords -M | mail -s uprecords $(whoami)
