# .config/systemd/system/athome.target
[Unit]
Description=dfanjul athome
Wants=boinc-client.service
Wants=cron.service
